# systemd/titan-commando-selector.service
[Unit]
Description=Titan NeuForge Commando Selector — v3.0 Full Lifecycle Governance
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/titan
EnvironmentFile=/opt/titan/.env
ExecStart=/usr/bin/python3 /opt/titan/Titan_CommandoSelector.py
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=titan-commando-selector
KillMode=process
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
