# systemd/titan-nf-outreach.service
[Unit]
Description=Titan NeuForge Outreach — 3-Sequence Email Engine (cold/FOMO/enterprise)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/titan
EnvironmentFile=/opt/titan/.env
ExecStart=/usr/bin/python3 /opt/titan/Titan_NeuForgeOutreach.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=titan-nf-outreach
KillMode=process
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
